# overlay/etc/init.d/ccapply
#!/sbin/openrc-run

depend() {
    want network-online
}

name="ccapply"
command="/bhojpur/system/os/current/opsutl"
command_args="config"